# ansible/templates/qlds@.service.j2
[Unit]
Description=Quake Live Dedicated Server Port %i
After=network.target

[Service]
User=ql
Group=ql
WorkingDirectory=/home/ql/qlds-%i
{% if cpu_affinity is defined and cpu_affinity is not none %}
CPUAffinity={{ cpu_affinity }}
{% endif %}
{% if ld_preload_paths | default('') %}
Environment=LD_PRELOAD={{ ld_preload_paths }}
{% endif %}
# The ExecStart line will include the script and the dynamic arguments passed from Ansible
# Wrap the entire command in quotes to handle spaces/quotes in qlds_args
ExecStart=/home/ql/qlds-%i/run_server_x64_minqlx.sh {{ qlds_args | default('') }} 
Restart=on-failure
RestartSec=5
# qzeroded/minqlx does not exit cleanly on SIGTERM. Treat fast hard-kill
# restarts as intentional for QLDS config/hook applies instead of burning a
# 10s stop timeout. Do not mark SIGKILL as a clean exit globally: external
# kill -9/OOM exits should still be failures under Restart=on-failure.
KillSignal=SIGKILL
TimeoutStopSec=2

# Standard output and error logging (optional, defaults to journald)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
